# openwrt/common/etc/rc.local
## 

/opt/lib/scripts/sg-data-mount.sh &
/opt/lib/scripts/sg-data-dirs.sh &

exit 0
