# openwrt/common/etc/rc.local
## 

## fix missing /dev/fd
[ -h /dev/fd ] || ln -s /proc/self/fd /dev/fd

## check SG specific dirs and mounts
/opt/lib/scripts/sg-data-mount.sh &
/opt/lib/scripts/sg-data-dirs.sh &

## check if wifi wan is up
/opt/lib/scripts/wifiwan-check.sh &

exit 0
