# openwrt/dir505/etc/rc.local
# some checks specific to SuperGlue

## if /www is mounted
mountpoint -q /www && (

  ## if not writable chown with httpd
  [ $(stat /www -c %U) != 'httpd' ] && chown -R httpd /www
  [ $(stat /www -c %a) -lt '755' ] && chmod -R u+rwX /www

  ## if /www/tmp does not exist
  [ ! -d /www/tmp ] && (
    mkdir /www/tmp
    chown -R httpd /www/tmp
  )
  ## if /www/htdocs does not exits
  [ ! -d /www/htdocs ] && (
    mkdir /www/htdocs
    chown -R httpd /www/htdocs
  )
)

exit 0
