# openwrt/common/etc/uci-defaults/21-gen-psk
/usr/sbin/openvpn --genkey --secret /opt/lib/host.psk 
