# openwrt/common/etc/uci-defaults/21-gen-psk
[ -f /etc/openvpn/host.psk ] || /usr/sbin/openvpn --genkey --secret /etc/openvpn/host.psk
