# openwrt/common/etc/uci-defaults/22-gen-htpasswd
[ -f /etc/lighttpd/htpasswd ] || (
  _PWD='changeme'
  _HASH=$(printf '%s' "admin:superglue:$_PWD" | md5sum | cut -d' ' -f1)
  echo "admin:superglue:$_HASH" > /etc/lighttpd/htpasswd
  echo -e "$_PWD\n$_PWD" | passwd root &>/dev/null
)
