# openwrt-rescue/common/etc/crontabs/root
1 6 * * *	  /sbin/logrotate /www/log/*.log; /etc/init.d/lighttpd reload
