# openwrt/DIR505A1/etc/uci-defaults/12-fstab-enable
[ -e /etc/config/fstab ] && [ -z $(find /etc/rc.d/ -name 'S??fstab') ] && ( 
  /etc/init.d/fstab enable
  /etc/init.d/fstab start
  )

