# openwrt/DIR505A1/etc/uci-defaults/20-set-ssid
_MAC=$(iw wlan0 info | awk 'BEGIN { FS=":" } /addr/ { print $5$6 }')                 
_MAC="SuperGlue-"$_MAC                                                               

[ $(uci get wireless.@wifi-iface[0].ssid) != $_MAC ] && (
  uci set wireless.@wifi-iface[0].ssid=$_MAC
  uci commit
  wifi
  /etc/init.d/network restart
  /etc/init.d/odhcpd restart
  )

exit 0
