# openwrt/TLWR710/etc/uci-defaults/20-set-ssid
_MAC=$(iw wlan0 info | awk 'BEGIN { FS=":" } /addr/ { print $5$6 }')                 
_MAC="SuperGlue-"$_MAC                                                               

[ $(uci get wireless.@wifi-iface[0].ssid) != $_MAC ] && (
  uci set wireless.@wifi-iface[0].ssid=$_MAC
  uci commit
  ( sleep 1; reboot; )&
  )

exit 0
