# openwrt/common/etc/crontabs/root
1 6 * * *	/root/logrotate.sh /www/htdocs/logs/access.log; /root/logrotate.sh /www/htdocs/logs/error.log; /etc/init.d/lighttpd restart

