# openwrt/common/etc/crontabs/root
## logrotation
1 6 * * *     /sbin/logrotate /www/log/*.log; /etc/init.d/lighttpd reload

## dyndns
*/5 * * * *   /opt/lib/scripts/dyndns-update.sh
