# openwrt/common/etc/crontabs/root
1 6 * * *	  /sbin/logrotate /www/htdocs/logs/access.log \
                            /www/htdocs/logs/error.log \
                            /www/htdocs/logs/post.log \
                            /www/htdocs/logs/admin.log; \
                            /etc/init.d/lighttpd restart
