# openwrt-rescue/common/etc/uci-defaults/20-set-ssid
_MAC=$(iw wlan0 info | awk 'BEGIN { FS=":" } /addr/ { print $5$6 }')                 
_MAC="Superglue-$_MAC-rescue"

[ $(uci get wireless.@wifi-iface[0].ssid) != $_MAC ] && (
  uci set wireless.@wifi-iface[0].ssid=$_MAC
  uci commit wireless
  wifi
  /etc/init.d/network restart
  /etc/init.d/dnsmasq restart
  )

exit 0
